# clamav-mirror/crontab
#Keep definitions up to date 
*/2 * * * * cvd update -V
#Test and publish signatures
*/50 * * * * ./test.sh